# slave/etc/systemd/system/bind9.service.d/override.conf
[Service]
ExecStartPost=/bin/sleep 5
ExecStartPost=/usr/sbin/rndc reload catalog.example.com
